# configs/99-broadcom-wl-openfyde.rules
ACTION=="add|change", SUBSYSTEM=="net", KERNEL=="wlan*", RUN+="/usr/local/sbin/broadcom-wl-openfyde-post-up %k"
