# configs/broadcom-wl-openfyde-autoconnect.conf
description "Broadcom wl WEXT Wi-Fi autoconnect"

start on (started shill and started ui)
stop on stopping shill

task
console log

exec /usr/local/sbin/broadcom-wl-openfyde-autoconnect
